# cron/crontab
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin
PYTHONPATH=/app
TZ=UTC
DB_HOST=mysql
DB_PORT=3306
DB_USER=root
DB_PASSWORD=root
DB_NAME=crypto_ai

0 3 * * * cd /app && flock -n /tmp/coingecko.lock python -m ai.src.cache_coingecko >> /app/logs/coingecko.log 2>&1
5 * * * * cd /app && flock -n /tmp/fetch_1h.lock python -m ai.src.fetch_1h_top300 >> /app/logs/fetch_1h.log 2>&1
0 4 * * * cd /app && flock -n /tmp/fetch_1d.lock python -m ai.src.fetch_1d_top300 >> /app/logs/fetch_1d.log 2>&1
0 6 * * * cd /app && flock -n /tmp/fetch_1w.lock python -m ai.src.fetch_weekly_all >> /app/logs/fetch_1w.log 2>&1
0 5 * * * cd /app && flock -n /tmp/train.lock python -m ai.src.train_all_auto 3 >> /app/logs/train.log 2>&1
*/10 * * * * cd /app && flock -n /tmp/predict.lock python -m ai.src.predict_all 3 >> /app/logs/predict.log 2>&1
*/10 * * * * cd /app && flock -n /tmp/evaluate.lock python -m ai.src.batch_evaluate >> /app/logs/evaluate.log 2>&1
*/30 * * * * cd /app && flock -n /tmp/overview_1h.lock python -m ai.jobs.build_market_overview 1h >> /app/logs/overview_1h.log 2>&1
15 5 * * * cd /app && flock -n /tmp/overview_1d.lock python -m ai.jobs.build_market_overview 1d >> /app/logs/overview_1d.log 2>&1
30 6 * * * cd /app && flock -n /tmp/overview_1w.lock python -m ai.jobs.build_market_overview 1w >> /app/logs/overview_1w.log 2>&1
*/5 * * * * echo "CRON OK $(date -u)" >> /app/logs/cron_health.log
